# Monitoring/ProcMon/procmon.service
[Unit]
Description=ProcMon

[Service]
TimeoutStartSec=0
ExecStart=/opt/monitoring/procmon/procmon

[Install]
WantedBy=multi-user.target