# Monitoring/DockMon/dockmon_supercronic.cron
# Run every minute
*/15 * * * * python3 ./dockmon.py