# Monitoring/DockMon/dockmon_supercronic.cron
# Run every 15 seconds
*/15 * * * * * python3 ./dockmon.py